# deploy/timekeeper-weekly.timer
[Unit]
Description=Send TimeKeeper weekly hours email

[Timer]
# Every Sunday at 8:00 PM, summarizing the week just ending (the service runs
# weekly_email.py --this-week). `Persistent` catches up if the Pi was off.
OnCalendar=Sun *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
